# eero-dashboard.service
[Unit]
Description=Eero Network Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/pi/eero-dashboard
Environment=PATH=/home/pi/eero-dashboard/venv/bin
ExecStart=/home/pi/eero-dashboard/venv/bin/python /home/pi/eero-dashboard/dashboard.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings (relaxed for port 80 access)
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/.eero-dashboard
ReadWritePaths=/root/.eero-dashboard

[Install]
WantedBy=multi-user.target